# archive/zfs-protection/systemd/zfs-scrub-scheduler.service
[Unit]
Description=ZFS Scrub Scheduler
Documentation=file:///opt/zfs-protection/docs/
After=zfs-import.target zfs-mount.service
Wants=zfs-import.target
PartOf=zfs.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/zfs-protection/scripts/zfs-scrub-manager.sh auto
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zfs-scrub-scheduler
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security settings
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true

# Resource limits
MemoryLimit=512M
CPUQuota=25%
TasksMax=50

# Timeout settings
TimeoutStartSec=300
TimeoutStopSec=60

# File system access
ReadWritePaths=/var/log/zfs-protection /var/run /tmp
ReadOnlyPaths=/etc/zfs-protection /opt/zfs-protection
PrivateTmp=true

[Install]
WantedBy=multi-user.target