# scripts/nfs-fileserver7-mount.service
# Systemd mount unit for fileserver7-nfs on FGSRV7
# Aponta para CT546 (fileserver7) em 192.168.70.240 — NÃO para fileserver5.
# Deploy: copiar para FGSRV7 como /etc/systemd/system/mnt-pve-fileserver7\x2dnfs.mount
#
# Também em /etc/fstab no FGSRV7:
# 192.168.70.240:/storage/nfs-export /mnt/pve/fileserver7-nfs nfs vers=3,hard,timeo=600,retrans=2,_netdev,nofail 0 0
# 192.168.70.240:/storage/fg_antigo  /mnt/fg-legacy-nfs       nfs vers=3,hard,timeo=600,retrans=2,_netdev,nofail 0 0

[Unit]
Description=Mount fileserver7 (CT546) NFS for Proxmox
Before=pve-cluster.service
After=network-online.target
Wants=network-online.target

[Mount]
What=192.168.70.240:/storage/nfs-export
Where=/mnt/pve/fileserver7-nfs
Type=nfs
Options=vers=3,hard,timeo=600,retrans=2,_netdev,nofail

[Install]
WantedBy=multi-user.target
