# config/systemd/makemoney01-ct/makemoney-crm-beleza.service
[Unit]
Description=makemoney01 crm-beleza (mvp)
After=network.target makemoney-lead-reply-api.service

[Service]
Type=simple
WorkingDirectory=/mnt/overpower/apps/dev/agl/crm-beleza/src
Environment=MAKEMONEY_LEAD_REPLY_URL=http://127.0.0.1:8765
Environment=MAKEMONEY_ROOT=/mnt/overpower/apps/dev/agl/makemoney01
Environment=DB_CONNECTION=sqlite
Environment=DB_DATABASE=/var/lib/makemoney/crm-beleza/database.sqlite
Environment=LARAVEL_STORAGE_PATH=/var/lib/makemoney/crm-beleza/storage
Environment=TMPDIR=/var/lib/makemoney/crm-beleza/tmp
ExecStart=/usr/bin/php8.4 artisan serve --host=127.0.0.1 --port=8106 --no-reload
Restart=on-failure
RestartSec=5
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
