# config/systemd/makemoney01-ct/makemoney-erp-estacionamento.service
[Unit]
Description=makemoney01 erp-estacionamento (mvp)
After=network.target makemoney-lead-reply-api.service

[Service]
Type=simple
WorkingDirectory=/mnt/overpower/apps/dev/agl/erp-estacionamento/src
Environment=MAKEMONEY_LEAD_REPLY_URL=http://127.0.0.1:8765
Environment=MAKEMONEY_ROOT=/mnt/overpower/apps/dev/agl/makemoney01
Environment=DB_CONNECTION=sqlite
Environment=DB_DATABASE=/var/lib/makemoney/erp-estacionamento/database.sqlite
Environment=LARAVEL_STORAGE_PATH=/var/lib/makemoney/erp-estacionamento/storage
Environment=TMPDIR=/var/lib/makemoney/erp-estacionamento/tmp
ExecStart=/usr/bin/php8.4 artisan serve --host=127.0.0.1 --port=8108 --no-reload
Restart=on-failure
RestartSec=5
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
