# config/systemd/makemoney01-ct/makemoney-crm-academia.service
[Unit]
Description=makemoney01 crm-academia demo (PHP built-in)
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/overpower/apps/dev/agl/crm-academia/src/public
ExecStart=/usr/bin/php -S 127.0.0.1:8110 -t /mnt/overpower/apps/dev/agl/crm-academia/src/public
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
