# config/systemd/makemoney01-ct/makemoney-erp-restaurante.service
[Unit]
Description=makemoney01 erp-restaurante demo (PHP built-in)
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/overpower/apps/dev/agl/erp-restaurante/src/public
ExecStart=/usr/bin/php -S 127.0.0.1:8111 -t /mnt/overpower/apps/dev/agl/erp-restaurante/src/public
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
