# config/systemd/makemoney01-ct/makemoney-erp-pme-finance.service
[Unit]
Description=makemoney01 erp-pme-finance (preview static)
After=network.target

[Service]
Type=simple
WorkingDirectory=/tmp/tmp5gtspfjm/erp-pme-finance
ExecStart=/usr/bin/php8.4 -S 127.0.0.1:8107 -t src/public
Restart=on-failure
RestartSec=5
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
